# etc/systemd/system/transmission.service.d/resources.conf
[Service]
CPUWeight=idle
IOSchedulingClass=idle
